# images/cron/crontab
* * * * *  php-console content:set "$(date)"
